# services/backup/run
#!/usr/bin/with-contenv /bin/bash
echo "Starting backup process"
exec python3 /backup.py
